# prog/dnsftpsmtp/indices/crontab
*/1 * * * * bash /var/tmp/indice1.sh  
*/1 * * * * bash /var/tmp/indice2.sh  
